# scripts/iso/files/systemd/openpalm-stack.timer
[Unit]
Description=Periodic OpenPalm stack reconciliation

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Unit=openpalm-stack.service

[Install]
WantedBy=timers.target
